# unison251/unison.service
[Unit]
Description=Unison service

[Service]
Type=exec
ExecStart=/usr/bin/unison -logfile=.unison/service.log -repeat=watch -auto -batch -terse

